# ops/systemd/repobrief-fleet/rb-publish-fleet-watch.service
[Unit]
Description=Publish changed RepoBrief fleet bundles with bounded retention

[Service]
Type=oneshot
ExecStart=/home/alex/.local/bin/rb-publish-fleet --if-changed --retention 3 --limit 8
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
